# examples/systemd/graphicalmatrix-csv-import.path
[Unit]
Description=Watch GraphicalMatrix provisioning CSV directory

[Path]
PathChanged=/opt/graphicalmatrix-admin/incoming
Unit=graphicalmatrix-csv-import.service

[Install]
WantedBy=multi-user.target
